# system/systemd-user/niri-power-effects.service
[Unit]
Description=Toggle Niri visual effects with the system power profile
PartOf=graphical-session.target
After=graphical-session.target power-profiles-daemon.service
Requisite=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.config/niri/scripts/power-profile-effects.sh
Restart=on-failure
RestartSec=3
