# system/systemd-user/niri-power-effects.path
[Unit]
Description=Watch the main Niri config for power-effect regeneration
PartOf=graphical-session.target

[Path]
PathChanged=%h/.config/niri/config.kdl
PathChanged=%h/.config/niri/conf.d/binds.kdl
Unit=niri-power-effects-refresh.service
