# system/systemd-user/niri-low-battery.service
[Unit]
Description=Switch to power saver on low battery
After=graphical-session.target power-profiles-daemon.service upower.service

[Service]
Type=simple
ExecStart=%h/.config/niri/scripts/low-battery-profile.sh
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
